# Semester_2/SisOp_Praktikum_1/challenge_1/crontabs
SHELL=/bin/bash

*/1 * * * * /bin/bash /home/blackcurrent/challange_1/fetch.sh >> /home/blackcurrent/challange_1/kalender.txt 2>&1
*/1 * * * * /bin/bash /home/blackcurrent/challange_1/miss.sh >> /home/blackcurrent/challange_1/kalender.txt 2>&1
